# scripts/password-safe-server.conf
description "Password safe server"

start on startup

task
exec /home/stephen/passwordSafeServer/dist/start  